# conf/s_naga.conf
[program:naga]
command=uwsgi --ini /home/www-data/Naga/uwsgi.ini
directory=/home/www-data/Naga
autostat=true
autorestat=true
