# conf/supervisord.conf
[program:naga]
command=uwsgi --ini /root/Naga/uwsgi.ini
directory=/root/naga
autostat=true
autorestat=true